# templates/systemd.timer.j2
[Unit]
Description=Let's Encrypt renewal

[Timer]
OnCalendar=daily
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
